# root/etc/services.d/bb-ws/run
#!/bin/bash

exec /usr/bin/php /app/barcodebuddy/wsserver.php